# others/warsaw/warsaw.service
[Unit]
Description=Handle Warsaw

[Service]
ExecStart=/usr/local/bin/warsaw/core

[Install]
WantedBy=multi-user.target
